# backend/syncsonic-runtime-sync.service
[Unit]
Description=SyncSonic Runtime Ultrasonic Measurement Service
After=syncsonic.service
Requires=syncsonic.service

[Service]
Type=simple
User=syncsonic
WorkingDirectory=/home/syncsonic/SyncSonicPi/backend
Environment=XDG_RUNTIME_DIR=/run/syncsonic
Environment=PULSE_SERVER=unix:/run/syncsonic/pulse/native
Environment=PULSE_SYSTEM_BUS=1
EnvironmentFile=-/etc/default/syncsonic
ExecStart=/usr/bin/env PYTHONUNBUFFERED=1 /home/syncsonic/SyncSonicPi/backend/.venv/bin/python -m syncsonic_ble.runtime_sync_service
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
